# my-elisp-files/tasks-web/com.kaan.tasks-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kaan.tasks-web</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/kaan/emacsdotfiles/my-elisp-files/tasks-web/server.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash; clean `launchctl unload' will still stop it. -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/tasks-web.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tasks-web.log</string>

    <!-- Override env vars here if needed.
    <key>EnvironmentVariables</key>
    <dict>
        <key>OBSIDIAN_DIR</key>
        <string>/Users/kaan/…</string>
        <key>TASKS_PORT</key>
        <string>8765</string>
    </dict>
    -->
</dict>
</plist>
